# packaging/deb/prerm
#!/bin/sh
set -e

echo "Stopping and disabling mprisence service..."
systemctl --user disable --now mprisence.service >/dev/null 2>&1 || :

exit 0
